# ops/launchd/com.floyd.surface.desktop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.floyd.surface.desktop</string>
  <key>ProgramArguments</key><array><string>/usr/local/bin/node</string><string>dist-server/index.js</string></array>
  <key>WorkingDirectory</key><string>/Volumes/Storage/FLOYD_WORKSTATION/intake/surfaces/desktop</string>
  <key>EnvironmentVariables</key><dict>
    <key>HOST</key><string>127.0.0.1</string><key>PORT</key><string>13010</string><key>MCP_WS_PORT</key><string>13011</string>
    <key>FLOYD_CORE_URL</key><string>http://127.0.0.1:41414</string><key>FLOYD_RUNTIME_ROOT</key><string>/Volumes/Storage/FLOYD_RUNTIME</string>
    <key>FLOYD_SURFACE_COMMIT</key><string>unverified</string>
  </dict>
  <key>RunAtLoad</key><true/><key>KeepAlive</key><true/><key>ThrottleInterval</key><integer>3</integer>
  <key>StandardOutPath</key><string>/Users/douglastalley/Library/Logs/Floyd/desktop.log</string>
  <key>StandardErrorPath</key><string>/Users/douglastalley/Library/Logs/Floyd/desktop.log</string>
</dict></plist>
