# com.floyd.core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.floyd.core</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Volumes/Storage/FLOYD_RUNTIME/releases/core/current/core/daemon/src/main.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Volumes/Storage/FLOYD_RUNTIME/releases/core/current</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/douglastalley</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>FLOYD_REMOTE_ORIGIN</key>
    <string>https://floyd-workstation-remote.captainphantasy.chatgpt.site</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>Umask</key>
  <integer>63</integer>
  <key>StandardOutPath</key>
  <string>/Users/douglastalley/Library/Logs/floyd/core.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/douglastalley/Library/Logs/floyd/core.err.log</string>
</dict>
</plist>
