# apps/frame/com.floyd.frame.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.floyd.frame</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Volumes/Storage/FLOYD_WORKSTATION/apps/frame/server/frame-server.mjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/Storage/FLOYD_WORKSTATION/apps/frame</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FRAME_PORT</key>
        <string>13030</string>
        <key>HOME</key>
        <string>/Users/douglastalley</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>Umask</key>
    <integer>63</integer>
    <key>StandardOutPath</key>
    <string>/Users/douglastalley/Library/Logs/Floyd/frame.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/douglastalley/Library/Logs/Floyd/frame.err.log</string>
</dict>
</plist>
